# src/mariadb/services/mysql-autobackup.service
[Unit]
Description=MySQL AutoBackup
After=mariadb.service
Requires=mariadb.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mysql-autobackup
ExecStop=/usr/local/bin/mysql-autobackup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target